# rootfs/etc/crontabs/nobody

* * * * * /usr/bin/php84  /var/www/html/admin/cli/cron.php > /dev/null