# config/etc/crontabs/nobody

* * * * * /usr/bin/php  /var/www/html/admin/cli/cron.php > /dev/null