# config/etc/crontabs/nobody
* * * * * echo "hello stackoverflow"$(date) >> /var/www/html/test_file 2>&1
* * * * * /usr/bin/php  /var/www/html/admin/cli/cron.php > /dev/null